# systemd/thirdcupd.service
[Unit]
Description=thirdcupd backend and system watchdog
Documentation=https://github.com/neiroun/thirdcupd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=thirdcupd
Group=thirdcupd
ExecStart=/usr/local/bin/thirdcupd -config /etc/thirdcupd/thirdcupd.json
Restart=always
RestartSec=5s
LogsDirectory=thirdcupd
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
